# supervisord.conf
[supervisord]
nodaemon=true

[program:backend]
command=./main

[program:react]
command=serve -p 3000 -s build